# docs/web/docker/supervisor/dev/app.conf
[program:front]
command=yarn dev
autostart=true
autorestart=true
priority=30
